# systemd/sandbox-manager@.service
[Unit]
Description=sandbox-manager (slot %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# %i is the REST port (9000 or 9002). Noise port = REST + 10 (9010 / 9012).
Environment=BIND_ADDR=127.0.0.1:%i
EnvironmentFile=/var/www/todoforai/apps/sandbox-manager/shared/.env
EnvironmentFile=/var/www/todoforai/apps/sandbox-manager/shared/noise.env
ExecStartPre=/bin/sh -c 'echo "NOISE_BIND_ADDR=127.0.0.1:$((%i + 10))" > /run/sandbox-manager-%i.noise'
EnvironmentFile=-/run/sandbox-manager-%i.noise
WorkingDirectory=/var/www/todoforai/apps/sandbox-manager/current
ExecStart=/var/www/todoforai/apps/sandbox-manager/current/sandbox-manager
Restart=on-failure
RestartSec=2
LimitNOFILE=65535

# Firecracker needs privileged access; run as root on the VM host.
# (Refuse to run this service on hosts that don't have /dev/kvm unless ENABLE_KVM=false in shared/.env)
User=root

[Install]
WantedBy=multi-user.target
